# scripts/battery_eco/systemd/battery_eco.service
#!/bin/sh

[Unit]
Description=My script
 
[Service]
ExecStart= /home/maxlttr/.config/scripts/battery_eco.sh
 
[Install]
WantedBy=multi-user.target

